# scripts/sci-assist-daily.service
[Unit]
Description=SCI Assist Bot Daily Message
Requires=sci-assist-bot.service

[Service]
Type=oneshot
User=aaron
Group=aaron
WorkingDirectory=/opt/sci-assist
EnvironmentFile=/opt/sci-assist/.env
ExecStart=/opt/sci-assist/scripts/daily_message.fish post
StandardOutput=journal
StandardError=journal
